# launchd/com.user.nix-gc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.nix-gc</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mads/dotfiles/scripts/nix-gc.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>14400</integer>

    <key>StandardOutPath</key>
    <string>/tmp/nix-gc-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/nix-gc-launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>